# service/data/ready-set.service.in
[Unit]
Description=ReadySet internal system service

[Service]
Type=dbus
BusName=@SERVICE_NAME@
User=@USER@
ExecStart=@LIBEXECDIR@/@SERVICE_NAME@
