# data/startup/ximper/setup.target
[Unit]
Description=OEM Setup Target
Requires=ready-set-ximper.service
After=multi-user.target

[Install]
WantedBy=default.target
