# data/ready-set-existing-user.service.in
[Unit]
Description=ReadySet setup fpr existing user

BindsTo=gnome-session.target
After=gnome-session.target

# Never run in GDM
ConditionUser=!@system
ConditionEnvironment=XDG_SESSION_CLASS=user

[Service]
Type=oneshot
ExecStart=@LIBEXECDIR@/ready-set
Restart=no
